# build/config/upstart.conf
# syslog-ng

description "an upstart task to make sure that syslog-ng is always running"
author "Richard Lee <lifuzu@gmail.com>"

start on startup
stop on shutdown

respawn

script
  service syslog-ng start
end script

